# supervisord-rsyslogd.conf
[program:rsyslogd]
command=/usr/sbin/rsyslogd -c5
numprocs=1
autostart=true
autorestart=true
